# etc/crontabs/root
*/15 * * * *	/usr/sbin/chs-db backup
*/1 * * * *	/usr/sbin/chs-cron watchdog
