# usr/lib/opkg/info/chsd.list
/etc/hotplug.d/iface/50-chs-wifidn
/usr/sbin/chs-cron
/etc/lighttpd/ssl/server.pem
/usr/sbin/chs-init
/etc/uci-defaults/70_chs-db-export-patch
/usr/sbin/chs-wifisetup
/etc/lighttpd/lighttpd-chsd.conf
/usr/sbin/chsd
/etc/uci-defaults/71_chs-rb751-migration
/etc/init.d/chs-boot
/etc/init.d/chs-fixtime
/usr/sbin/chs-wifirestore
/root/.chsd/chsd.sql
/etc/init.d/chsd
/usr/sbin/chs-wifiscan
/usr/sbin/chs-version
/etc/config/chsd
/usr/sbin/chs-db
/usr/sbin/chs-update
/etc/avahi/services/clage-hs.service
